# weio/fs-overlay/etc/systemd/system/wifi_ap.service
[Unit]
Description=AccessPoint daemon
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/hostapd /etc/hostapd.conf

[Install]
WantedBy=multi-user.target
